# systemd/signal-claude-bridge.service
[Unit]
Description=Signal-to-Claude-Code Bridge
Wants=network-online.target
After=network-online.target signal-cli-daemon.service
Requires=signal-cli-daemon.service

[Service]
SyslogIdentifier=signal-claude-bridge
Type=simple
# Adjust WorkingDirectory and ExecStart to your installation path
WorkingDirectory=/opt/signal-claude-bridge
ExecStart=/usr/bin/python3 /opt/signal-claude-bridge/signal_claude_bridge.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
# Ensure claude CLI is in PATH
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=%h"

# Security hardening — restrict what the bridge process can access
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
# Claude's Bun runtime and Java both require JIT compilation
MemoryDenyWriteExecute=no
# Adjust these paths to match your deployment:
# - Claude subprocess needs read-write access to its data directory
# - Add any directories your agent needs to write to
ReadWritePaths=%h/.claude
# - Config and bridge settings should be read-only at runtime
ReadOnlyPaths=%h/.config/system-config.json

[Install]
WantedBy=multi-user.target
